# deploy/systemd/newsbeat-digest.service
[Unit]
Description=Run the newsbeat-digest batch processor
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/newsbeat-digest
ExecStart=/usr/bin/flock -n /run/lock/newsbeat-digest.lock /usr/bin/docker compose run --rm newsbeat-digest
TimeoutStartSec=45min
Nice=10
IOSchedulingClass=idle
