# deploy/systemd/newsbeat-digest.timer
[Unit]
Description=Run newsbeat-digest four times daily

[Timer]
OnCalendar=*-*-* 07:07:00 Europe/Tallinn
OnCalendar=*-*-* 11:07:00 Europe/Tallinn
OnCalendar=*-*-* 16:07:00 Europe/Tallinn
OnCalendar=*-*-* 21:07:00 Europe/Tallinn
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
